# config/launchd/com.bopvault.ai-research-assistant.clips-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bopvault.ai-research-assistant.clips-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/boris.diebold/code/ai-research-assistant/scripts/clips-watcher.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/boris.diebold/code/ai-research-assistant/logs/clips-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/boris.diebold/code/ai-research-assistant/logs/clips-watcher-error.log</string>
</dict>
</plist>
